# packages/cron/etc/services.d/cron/run
#!/usr/bin/with-contenv bash

cron_bin=$(command -v cron)
if [[ -z ${cron_bin} ]]; then
    cron_bin=$(command -v crond)
fi

exec ${cron_bin} -f
